# conf.d/oleg_test8.conf
[program:oleg_test8]
command = python3 bots/oleg_test8.py
environment = PYTHONUNBUFFERED=1
autostart = false
stdout_logfile = logs/oleg_test8.log
stderr_logfile = logs/oleg_test8_error.log

